# systemd/premarket-overnight-update.service
[Unit]
Description=OpenClaw US premarket overnight update

[Service]
Type=oneshot
WorkingDirectory=/root/.openclaw/workspace/us-premarket-telegram-reports
Environment=TELEGRAM_TARGET=8651624513
ExecStart=/root/.openclaw/workspace/us-premarket-telegram-reports/scripts/premarket_overnight_update.sh
